# packages/cip-home/services/marsbar.service
[Install]
WantedBy=xsession.target

[Service]
ExecStart=marsbar

[Unit]
Description=Bar for marswm desktop
