# packages/cip-home/services/wallpaper-daemon.service
[Install]
WantedBy=xsession.target

[Service]
ExecStart=wallpaper-daemon

[Unit]
Description=Watch for configuration changes in the display
